# com.roblox.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>             <string>com.roblox.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/var/jb/usr/local/bin/roblox-watchdog.sh</string>
    </array>

    <!-- Run at boot and every 60 s -->
    <key>RunAtLoad</key>         <true/>
    <key>StartInterval</key>     <integer>60</integer>

    <!-- GUI user context -->
    <key>LimitLoadToSessionType</key> <string>Aqua</string>

    <!-- send output to unified logging – files no longer needed -->
    <key>StandardOutPath</key>   <string>/dev/null</string>
    <key>StandardErrorPath</key> <string>/dev/null</string>
</dict>
</plist>
