# com.roblox.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>            <string>com.roblox.watchdog</string>
    <key>ProgramArguments</key> <array>
        <string>/var/jb/usr/local/bin/roblox-watchdog.sh</string>
    </array>

    <!-- fire every 60 s and once right after bootstrap -->
    <key>StartInterval</key>    <integer>60</integer>
    <key>RunAtLoad</key>        <true/>

    <!-- drop privileges to “mobile” but remain a LaunchDaemon -->
    <key>UserName</key>         <string>mobile</string>

    <!-- stdout/stderr go through ASL/syslog; nothing here -->
</dict>
</plist>